# systemd/horizon-backend.service
[Unit]
Description=Horizon Institutional Lab backend workers
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/horizon-lab
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml --env-file .env up -d --build
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml --env-file .env down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
